# rootfs/etc/cron.daily/certbot
#!/usr/bin/env sh

exec /opt/certbot/certbot-auto renew --quiet --pre-hook /opt/certbot/pre-hook --post-hook /opt/certbot/post-hook --renew-hook /opt/certbot/renew-hook
