# rootfs/etc/cron.daily/certbot
#!/usr/bin/env sh

exec /opt/certbot/certbot-auto renew
