# rootfs/etc/cron.daily/certbot
#!/usr/bin/env sh

/opt/certbot/certbot-auto renew --quiet
