# rootfs/etc/cron.daily/duplicity
#!/usr/bin/env sh

/bootstrap.sh/scripts/cronic /bootstrap.sh/externals/duply/duply awsS3 backup
