# rootfs/etc/cron.weekly/duplicity
#!/usr/bin/env bash

export USER=root
source $HOME/.nix-profile/etc/profile.d/nix.sh

# duply-cronic awsS3 purgeIncr_purgeFull_purge --force
