# rootfs/etc/cron.weekly/duplicity
#!/usr/bin/env sh

# /bootstrap.sh/externals/duply/duply awsS3 purgeIncr_purgeFull_purge --force
